# etc/systemd/system/mnt-cluster.automount
[Unit]
Description=NFS automount for DL cluster

[Automount]
Where=/mnt/cluster

[Install]
WantedBy=multi-user.target
